# redis/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usb/sbin/redis-server /etc/redis/redis.conf
autorestart=true
